# service/haproxy/run
#!/bin/bash

/commands/certbot-concat-crt
exec haproxy -f /etc/haproxy/haproxy.cfg
